# devops/ansible/roles/web/templates/nginx_supervisor.conf.j2
[program:nginx]
command=service nginx start
autostart=true
autorestart=true
redirect_stderr=false
